# crontab
30 * * * * cd code/qbtprune && ./cronjob
