# crontab
0 4,10,16,22 * * * cd code/qbtprune && ./cronjob
